# babymonitor.service.in
[Unit]
Description=Baby Monitoring Service

[Service]
Type=simple
ExecStart=@bindir@/babymonitor --config @sysconfdir@/babymonitor/config.ini
CapabilityBoundingSet=
ProtectSystem=true
ProtectHome=false
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
